# containers/ftp-downloader/scripts/crontab
0 2 * * * /scripts/download.sh >> /var/log/ftp-download.log 2>&1